# etc/cron.d/jagc
*/1 * * * * root su jagc -c "cd /home/jagc/jagc.org && ./script/jagc stat -m production"
